# NetFavorite/饿了么/别人的代码/基准/hackathon-2015-master/hackathon-vagrant/files/hackathon-appreload.conf
respawn
respawn limit 5 5

start on vagrant-mounted
stop on runlevel [06]

limit nofile 65535 65535

exec /usr/local/bin/hackathon-appreload >> /vagrant/app.log 2>&1
